# etc/systemd/system/xinit-login.service
[Unit]
After=systemd-user-sessions.service

[Service]
ExecStart=/home/ldh/.local/bin/xinit-login.sh

[Install]
WantedBy=multi-user.target
